# templates/default/timer-weekly.timer.erb
[Unit]
Description=Weekly Timer

[Timer]
OnBootSec=15min
OnCalendar=weekly
Persistent=true
Unit=timer-weekly.target

[Install]
WantedBy=basic.target
