# templates/default/timer-yearly.timer.erb
[Unit]
Description=Yearly Timer

[Timer]
OnBootSec=90min
OnCalendar=yearly
Persistent=true
Unit=timer-yearly.target

[Install]
WantedBy=basic.target
