# templates/default/timer-yearly.target.erb
[Unit]
Description=Yearly Timer Target
StopWhenUnneeded=yes
