# templates/default/timer-daily.timer.erb
[Unit]
Description=Daily Timer

[Timer]
OnBootSec=10min
OnCalendar=daily
Persistent=true
Unit=timer-daily.target

[Install]
WantedBy=basic.target
