# templates/default/timer-monthly.timer.erb
[Unit]
Description=Monthly Timer

[Timer]
OnBootSec=45min
OnCalendar=monthly
Persistent=true
Unit=timer-monthly.target

[Install]
WantedBy=basic.target
