# templates/default/timer-monthly.target.erb
[Unit]
Description=Monthly Timer Target
StopWhenUnneeded=yes
